# basics/config/includes.chroot/etc/profile.d/disable-gnome-welcome.sh
#!/bin/sh

set -e

echo "yes" > ${HOME}/.config/gnome-initial-setup-done
